# ansible/roles/wayvnc/files/wayvnc.service
[Unit]
Description=WayVNC on headless Hyprland output
After=graphical-session.target
BindsTo=graphical-session.target

[Service]
Type=simple
ExecStartPre=/usr/bin/hyprctl output create headless VNC-1
ExecStartPre=/bin/sleep 0.2
ExecStartPre=/usr/bin/hyprctl keyword monitor "VNC-1,1920x1080@60,auto,1"
ExecStart=/bin/bash -c 'exec /usr/bin/wayvnc -o VNC-1 "$(/usr/bin/tailscale ip -4)" 5900'
ExecStopPost=/usr/bin/hyprctl output remove VNC-1
Restart=on-failure
RestartSec=3

[Install]
WantedBy=graphical-session.target
