# ansible/roles/maintenance/files/dotmachines-reconcile.service
[Unit]
Description=dotmachines config reconciliation
After=network-online.target
Wants=network-online.target
ConditionPathExists=/home/atqa/repo/dotmachines/ansible/playbooks/bootstrap.yaml

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/dotmachines-reconcile.sh
TimeoutStartSec=30min
Nice=10
IOSchedulingClass=idle

[Install]
WantedBy=multi-user.target
